# mbus-health-monitor.service
[Unit]
Description=MBus Gateway Health Monitor
After=mbus-mqtt-gateway.service
Requires=mbus-mqtt-gateway.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/mbus-mqtt-gateway
ExecStart=/opt/mbus-mqtt-gateway/.venv/bin/python health_monitor.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
